# beaglebone/os/board/wfe/beaglebone/rootfs_overlay/etc/crontab
@reboot python /home/wfe/app/wfe/can_monitor.py -d /home/wfe/data/2018CAR.dbc & 
